# _rootfs/etc/services.d/nginx/run
#!/bin/bash

if ! [[ -e /tmp/.ran_once ]]; then
	touch /tmp/.ran_once
	mkdir -p /var/log/nginx
fi

exec nginx
