# _rootfs/etc/services.d/nginx/run
#!/bin/bash

mkdir -p /var/log/nginx

sed -i "s|<USER> <GROUP>|etdev users|g;
		s|<WP_VERSION>|${WP_VERSION}|g" /etc/nginx/nginx.conf

exec nginx
